# ssh/crontab
# Logout ssh-agent every 8 hours
0 */8 * * * . ~/.ssh/agent-info && /usr/local/bin/ssh-agent-logout >> /var/log/ssh-agent.log 2>&1
