# mbsync/crontab
# Fetch mail every 5 minutes
*/5 * * * * /usr/local/bin/fetch-mail >> /var/log/fetch-mail.log 2>&1
